# scripts/old/sipcapture.crontab
##################################################################
# Running rotation script to cron once a day (scripts/homer_rotate)
# at low traffic time

30 4 * * * root /usr/bin/homer_rotate > /dev/null 2>&1
